# deploy/supervisor/fiddler_tools.conf
[program:fiddler_tools]
command=/usr/local/bin/gunicorn -w10 -b127.0.0.1:8888 flask_main_online:app
directory=/fiddler_tools/server
startsecs=5
stopwaitsecs=0
autostart=true
autorestart=true
